# systemd/ai-health-board.service
[Unit]
Description=AI Health Status Board
After=network.target

[Service]
Type=simple
# NOTE: This is a template. Paths and User are set by scripts/install_services.sh
User=__USER__
Group=__USER__
WorkingDirectory=__REPO_DIR__
Environment="PATH=__REPO_DIR__/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="GPIOZERO_PIN_FACTORY=lgpio"
ExecStart=__REPO_DIR__/venv/bin/python __REPO_DIR__/app.py run --config __REPO_DIR__/config/providers.yaml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
