# systemd/pisugar-button.service
[Unit]
Description=PiSugar S Button Daemon (GPIO3)
After=multi-user.target

[Service]
Type=simple
# NOTE: This is a template. Paths and User are set by scripts/install_services.sh
User=__USER__
Group=__USER__
ExecStart=/usr/bin/python3 __REPO_DIR__/scripts/pisugar_button.py
Environment=GPIOZERO_PIN_FACTORY=lgpio
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
